# znm2-libwrt-nss-no-wifi-usb/files/etc/uci-defaults/99-custom-defaults
#!/bin/sh
uci set network.lan.ipaddr='192.168.10.1'
uci set luci.main.lang='zh_cn'
uci commit network
uci commit luci
exit 0
